# snapshot_2025-08-20/filtered_json/train/cc90652d-048b-4e4d-923c-5ab8511ae51f.json
{
    "Event":{
        "threat_level_id":"2",
        "date":"2024-05-22",
        "published":true,
        "info":"Tiny BackDoor Goes Undetected – Suspected Turla leveraging MSBuild to Evade detection",
        "publish_timestamp":"1716360864",
        "timestamp":"1716360857",
        "Tag":[
            {
                "name":"misp-galaxy:mitre-attack-pattern=\"Deobfuscate/Decode Files or Information - T1140\""
            },
            {
                "name":"misp-galaxy:mitre-attack-pattern=\"Exploitation for Client Execution - T1203\""
            },
            {
                "name":"misp-galaxy:mitre-attack-pattern=\"Exfiltration Over C2 Channel - T1646\""
            },
            {
                "name":"misp-galaxy:mitre-attack-pattern=\"Application Layer Protocol - T1071\""
            },
            {
                "name":"misp-galaxy:mitre-attack-pattern=\"Scheduled Task - T1053.005\""
            },
            {
                "name":"misp-galaxy:mitre-attack-pattern=\"MSBuild - T1127.001\""
            },
            {
                "name":"misp-galaxy:mitre-attack-pattern=\"Masquerading - T1036\""
            },
            {
                "name":"tlp:white"
            }
        ],
        "Attribute":[
            {
                "category":"External analysis",
                "comment":"",
                "timestamp":"1716360638",
                "type":"link",
                "value":"https://cyble.com/blog/tiny-backdoor-goes-undetected-suspected-turla-leveraging-msbuild-to-evade-detection/",
                "to_ids":false
            },
            {
                "category":"Payload delivery",
                "comment":"",
                "timestamp":"1716360724",
                "type":"sha256",
                "value":"<sha256>",
                "to_ids":true
            },
            {
                "category":"Payload delivery",
                "comment":"",
                "timestamp":"1716360724",
                "type":"sha256",
                "value":"<sha256>",
                "to_ids":true
            },
            {
                "category":"Payload delivery",
                "comment":"",
                "timestamp":"1716360724",
                "type":"sha256",
                "value":"<sha256>",
                "to_ids":true
            },
            {
                "category":"Payload delivery",
                "comment":"",
                "timestamp":"1716360724",
                "type":"sha256",
                "value":"<sha256>",
                "to_ids":true
            },
            {
                "category":"Payload delivery",
                "comment":"",
                "timestamp":"1716360724",
                "type":"sha256",
                "value":"<sha256>",
                "to_ids":true
            },
            {
                "category":"Payload delivery",
                "comment":"",
                "timestamp":"1716360724",
                "type":"sha256",
                "value":"<sha256>",
                "to_ids":true
            }
        ]
    }
}