# frozen_snapshot/snapshot_2025-08-20/filtered_json/test/5fd62412-4cec-450a-a2bc-e52ec0a8ab16.json
{
    "Event":{
        "threat_level_id":"2",
        "date":"2020-12-13",
        "published":true,
        "info":"xHunt Campaign: Newly Discovered Backdoors Using Deleted Email Drafts and DNS Tunneling for Command and Control",
        "publish_timestamp":"1611866658",
        "timestamp":"1607869658",
        "Tag":[
            {
                "name":"tlp:white"
            },
            {
                "name":"misp-galaxy:mitre-enterprise-attack-attack-pattern=\"Scheduled Task - T1053\""
            }
        ],
        "Attribute":[
            {
                "category":"Payload delivery",
                "comment":"TriFive Samples",
                "timestamp":"1607869509",
                "type":"sha256",
                "value":"<sha256>",
                "to_ids":true
            },
            {
                "category":"Payload delivery",
                "comment":"Snugy Samples",
                "timestamp":"1607869509",
                "type":"sha256",
                "value":"<sha256>",
                "to_ids":true
            },
            {
                "category":"Payload delivery",
                "comment":"Snugy Samples",
                "timestamp":"1607869509",
                "type":"sha256",
                "value":"<sha256>",
                "to_ids":true
            },
            {
                "category":"Payload delivery",
                "comment":"Snugy Samples",
                "timestamp":"1607869509",
                "type":"sha256",
                "value":"<sha256>",
                "to_ids":true
            },
            {
                "category":"Payload delivery",
                "comment":"Snugy Samples",
                "timestamp":"1607869509",
                "type":"sha256",
                "value":"<sha256>",
                "to_ids":true
            },
            {
                "category":"Network activity",
                "comment":"Snugy C2 Domains",
                "timestamp":"1607869509",
                "type":"domain",
                "value":"deman1.icu",
                "to_ids":true
            },
            {
                "category":"Network activity",
                "comment":"Snugy C2 Domains",
                "timestamp":"1607869509",
                "type":"domain",
                "value":"hotsoft.icu",
                "to_ids":true
            },
            {
                "category":"Network activity",
                "comment":"Snugy C2 Domains",
                "timestamp":"1607869509",
                "type":"domain",
                "value":"uplearn.top",
                "to_ids":true
            },
            {
                "category":"Network activity",
                "comment":"Snugy C2 Domains",
                "timestamp":"1607869509",
                "type":"domain",
                "value":"lidarcc.icu",
                "to_ids":true
            },
            {
                "category":"Network activity",
                "comment":"Snugy C2 Domains",
                "timestamp":"1607869509",
                "type":"domain",
                "value":"sharepoint-web.com",
                "to_ids":true
            },
            {
                "category":"External analysis",
                "comment":"",
                "timestamp":"1607869524",
                "type":"url",
                "value":"https://unit42.paloaltonetworks.com/xhunt-campaign-backdoors/",
                "to_ids":false
            }
        ]
    }
}